# packaging/ghostty-voice.install
post_install() {
  cat <<'EOF'
==> ghostty-voice installed. Per-user setup (cannot run from package install):

    1. Copy and edit the config:
         install -Dm644 /usr/share/ghostty-voice/config.toml.example \
           ~/.config/ghostty-voice/config.toml

    2. Enable the daemon (it fetches the ~3 GB model on first run):
         systemctl --user enable --now ghostty-voiced
       For dictation before/without a graphical login:
         loginctl enable-linger "$USER"

    3. Drive your agent with the PTY wrapper — `talk-to` IS the interface.
       It wraps the agent on a pseudo-terminal and registers a wrapper sink, so
       your dictation is injected straight into that agent's PTY (no trailing
       Enter):
         talk-to claude
         talk-to ssh host claude
       While you are in the talk-to window: Shift+F10 toggles a batch recording,
       Shift+F9 starts a hands-free VAD recording. Triggers fire ONLY in that
       window — there is no system-wide hotkey. Cancel with
       `ghostty-voice-ctl cancel`.

    4. Diagnose (is the daemon reachable?):
         ghostty-voice-ctl doctor

    See /usr/share/doc/ghostty-voice/README.md for the full guide.
EOF
}

post_upgrade() {
  cat <<'EOF'
==> ghostty-voice upgraded.

    IMPORTANT — restart the daemon to pick up the new binary:

         systemctl --user restart ghostty-voiced

    A running ghostty-voiced keeps the OLD binary in memory until it is
    restarted. A stale daemon speaks an older control protocol, so `talk-to`
    shows `incompatible` and dictation will not reach the wrapped agent until
    you restart. Run the command above as the user whose daemon is running
    (not root).
EOF

  # Best-effort: restart any ALREADY-RUNNING per-user instances so the new binary
  # takes effect without manual action. Package scriptlets run as root and cannot
  # trivially reach a per-user systemd service, so this is strictly best-effort
  # and must never (a) fail the package transaction or (b) START a daemon that
  # isn't already running. We use:
  #   - `try-restart`, which acts ONLY on an already-active unit (no-op otherwise,
  #     so a disabled/stopped daemon is never started),
  #   - `--machine=<user>@.host` to reach each logged-in user's systemd manager,
  #   - and swallow all output/errors (|| true) so the transaction always succeeds.
  if command -v systemctl >/dev/null 2>&1 && command -v loginctl >/dev/null 2>&1; then
    local _gv_user
    while read -r _gv_user; do
      [ -n "$_gv_user" ] || continue
      systemctl --machine="${_gv_user}@.host" --user try-restart ghostty-voiced.service \
        >/dev/null 2>&1 || true
    done < <(loginctl list-users --no-legend 2>/dev/null | awk '{print $2}')
  fi
  return 0
}
